# recipes-support/usrp-hwd/files/usrp-hwd.service
[Unit]
Description=USRP Hardware Daemon (MPM)

[Service]
ExecStart=/usr/bin/usrp_hwd.py

[Install]
WantedBy=multi-user.target
