# meta-titanium/recipes-devtools/ledctrl/files/ledctrl-led1.service
[Unit]
Description=Configuring LED1 trigger

[Service]
Type=oneshot
ExecStart=/sbin/ledctrl led1 fpga

[Install]
WantedBy=multi-user.target
