# meta-titanium/recipes-devtools/ledctrl/files/ledctrl-led2.service
[Unit]
Description=Configuring LED2 trigger

[Service]
Type=oneshot
ExecStart=/sbin/ledctrl led2 emmc

[Install]
WantedBy=multi-user.target
