# meta-titanium/recipes-devtools/ledctrl/files/ledctrl-led0.service
[Unit]
Description=Configuring LED0 trigger

[Service]
Type=oneshot
ExecStart=/sbin/ledctrl led0 heartbeat

[Install]
WantedBy=multi-user.target
